# templates/open-webui.service.j2
[Unit]
Description=Open WebUI - Web interface for Ollama
After=network.target ollama.service

[Service]
Type=simple
User={{ ollama__open_webui_user }}
Group={{ ollama__open_webui_user }}
WorkingDirectory={{ open_webui_home }}
Environment="PATH={{ open_webui_venv_path }}/bin:/usr/local/bin:/usr/bin:/bin"
Environment="OLLAMA_BASE_URL=http://{{ ollama__host }}:{{ ollama__port }}"
Environment="WEBUI_AUTH={{ ollama__open_webui_auth | lower }}"
Environment="PORT={{ ollama__open_webui_port }}"
ExecStart={{ open_webui_venv_path }}/bin/open-webui serve
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
